# deploy/monit.service
[Unit]
Description=framebuffer cluster dashboard (local + remote GPU host)
After=network-online.target
Wants=network-online.target
# Take tty1 away from the login prompt while the dashboard owns the screen.
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
# Ensure the Super I/O sensor/fan driver is present for the thermal governor.
ExecStartPre=-/sbin/modprobe nct6775
ExecStart=/usr/local/bin/monit
# Site-specific settings live here (hostnames, ssh target). Not in the repo.
Environment=MONIT_CONFIG=/etc/monit/monit.conf
Restart=always
RestartSec=3
# Needs /dev/fb0 and /dev/tty0; runs as root.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
